# scripts/launchd/com.baker.forge-snapshot-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.baker.forge-snapshot-push</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__WORKER_PATH__</string>
  </array>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>RunAtLoad</key>
  <true/>
  <!-- BRISEN_LAB_FORGE_TELEMETRY_DURABILITY_1 (T3): auto-resurrect the pusher so
       it never dies silently (outage 2026-07-04). RunAtLoad resumes it on
       reboot/login; KeepAlive resumes it on crash. NOTE (flagged to lead): this
       job is a fast-exit periodic script (StartInterval=30). With KeepAlive=true
       launchd relaunches on EVERY exit, throttled to ~10s — effectively a ~10s
       cadence (≈3x the intended 30s) plus respawn log spam. If the 30s cadence
       matters, switch to the crash-only dict form which preserves StartInterval:
         <key>KeepAlive</key><dict><key>SuccessfulExit</key><false/></dict> -->
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/dimitry/Library/Logs/forge-snapshot-push.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dimitry/Library/Logs/forge-snapshot-push.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>FORGE_KEY</key>
    <string>__FORGE_KEY__</string>
    <key>LAB_URL</key>
    <string>https://brisen-lab.onrender.com</string>
  </dict>
</dict>
</plist>
