# launchd/com.brisen.forge-drift-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brisen.forge-drift-check</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__WRAPPER__</string>
    </array>
    <!-- Daily at 08:15 local. Not RunAtLoad — a boot-time run would race the
         session-start provisioning and false-alarm during login. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__LOGDIR__/forge-drift-check.stdout</string>
    <key>StandardErrorPath</key>
    <string>__LOGDIR__/forge-drift-check.stderr</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>FORGE_CHECK_DIR</key>
        <string>__BUNDLE__</string>
    </dict>
</dict>
</plist>
